# voice-transcriber.service
[Unit]
Description=Discord Voice Transcriber Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/mnt/data/voice-transcriber
ExecStart=/usr/bin/python3 /mnt/data/voice-transcriber/bot.py
EnvironmentFile=/mnt/data/voice-transcriber/.env
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
